# AI-Agent/ai_suricata_blocker.service
[Unit]
Description=AI Suricata Log Analyzer and IP Blocker
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/ai_suricata_blocker.py
Restart=always
User=root
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ai_suricata_blocker

[Install]
WantedBy=multi-user.target
